# ops/systemd/flipdot-bluetooth-tune.service
[Unit]
Description=Tune wireless adapters for game controller stability
After=bluetooth.service NetworkManager.service flipdot-bluetooth-ertm.service
Requires=bluetooth.service
Before=flipdot.service

[Service]
Type=oneshot
RemainAfterExit=yes
# Keep the adapter out of HID low-power link policies that can add latency or
# cause fragile controller links to churn. The IINE mini controllers identify as
# Xbox-compatible HID devices, but their Bluetooth firmware appears sensitive to
# the same low-power policy and retransmission behavior. This favors controller
# stability over tiny Bluetooth power savings.
ExecStart=-/usr/bin/hciconfig hci0 lp rswitch
# The Bluetooth and Wi-Fi radios share a Realtek combo chip. Keeping Wi-Fi power
# save disabled avoids coexistence stalls that can show up as BLE supervision
# timeouts on controller links.
ExecStart=-/usr/sbin/iw dev wlP1p1s0 set power_save off

[Install]
WantedBy=multi-user.target
